# ss_config/ss_supervisor.conf
[program:shadowsocks]
command=sslocal -c /etc/supervisor/conf.d/ss_bwg_2017.json
autostart=true
autorestart=true
user=root
log_stderr=true
logfile=/var/log/shadowsocks.log
